# opendap-app/supervisord.conf
[supervisord]
nodaemon=true

[program:opendap_server]
command= bash /entrypoint.sh

[program:new_thing]
command= python /root/app.py

[program: workers_cron]
command= python /root/cron_workers.py